# scripts/crontab
0 3 * * * /opt/alertsystem/scripts/alerts_check.sh
# Empty line at the end must be kept
