# scripts/crontab
0 3 * * * python3 /opt/alertsystem/manage.py check_alerts
# Empty line at the end must be kept
